# bootstrap/firstboot.service
[Unit]
Description=Initial network configuration

# machine-id is generated by the installer so "FirstBoot" doesn't happen
#ConditionFirstBoot=yes
#Before=first-boot-complete.target
#Wants=first-boot-complete.target

ConditionPathExists=/etc/do-run-firstboot

[Service]
ExecStart=/usr/local/sbin/firstboot.py

[Install]
WantedBy=multi-user.target
